# systemd/inkbird-ibs-p01r-mqtt.service
[Unit]
Description=Inkbird IBS-P01R SDR MQTT Decoder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=inkbird
Group=inkbird
RuntimeDirectory=inkbird-ibs-p01r
RuntimeDirectoryMode=0755
WorkingDirectory=/opt/inkbird-ibs-p01r-mqtt
EnvironmentFile=-/etc/inkbird-to-mqtt.env
ExecStart=/opt/inkbird-ibs-p01r-mqtt/.venv/bin/inkbird-ibs-p01r-mqtt run
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
